# autostart/startup.desktop
[Desktop Entry]
Exec=/home/excyber/script/startup.sh
