# autostart/startup.sh.desktop
[Desktop Entry]
Exec=/home/excyber/script/startup.sh
Icon=
Name=startup.sh
Path=
Terminal=False
Type=Application
