# config/autostart/xkeyboard.desktop
[Desktop Entry]
Hidden=true
Exec=xkeyboard
